# ansible-registry/roles/docker_registry/templates/registry-gc.timer.j2
[Unit]
Description=Run Docker Registry Garbage Collection

[Timer]
OnCalendar={{ gc_calendar }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
