# ansible-registry/roles/docker_registry/templates/registry-gc.service.j2
[Unit]
Description=Docker Registry Garbage Collection
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/registry-gc.sh
StandardOutput=journal
StandardError=journal
